# crontab
* * * * 1-5 cp /etc/hostsoff /etc/hosts

* * * * 6,0 cp /etc/hostson /etc/host
